# infra/supervisor/catalyst-worker.conf
[program:catalyst-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/catalyst/apps/api/artisan queue:work --memory=768 --tries=3 --max-time=3600
directory=/var/www/catalyst/apps/api
user=www-data
numprocs=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
; queue:work traps SIGTERM and finishes the in-flight job before exiting.
; Give it long enough for a worst-case 50 MP image decode:
stopwaitsecs=120
redirect_stderr=true
stdout_logfile=/var/log/catalyst/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
