# infra/supervisor/catalyst-scheduler.conf
[program:catalyst-scheduler]
process_name=%(program_name)s
command=php /var/www/catalyst/apps/api/artisan schedule:work
directory=/var/www/catalyst/apps/api
user=www-data
numprocs=1
autostart=true
autorestart=true
; schedule:work may be mid-command when told to stop; let a tick finish.
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=/var/log/catalyst/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
